# etc/systemd/system/nateflix.service
[Unit]
Description=Nateflix
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/docker-compose -f /home/nate/code/nateflix/docker-compose.yml up --remove-orphans -d
ExecStartPre=/usr/bin/docker-compose -f /home/nate/code/nateflix/docker-compose.yml pull
ExecStop=/usr/bin/docker-compose -f /home/nate/code/nateflix/docker-compose.yml down

[Install]
WantedBy=multi-user.target
